# nvmboot.service
[Service]
ExecStart=/usr/local/bin/tpm2_getsession

[Install]
WantedBy=multi-user.target


